# 2H2G3M/systemd/bs01-gunicorn.service
[Unit]
Description=BS01 Django Backend (Gunicorn) - 2C2G tuned
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/BS01/backend
Environment=PYTHONUNBUFFERED=1
Environment=DJANGO_SETTINGS_MODULE=backend.settings
EnvironmentFile=/root/BS01/backend/.env
ExecStart=/root/BS01/.venv/bin/gunicorn -c /root/BS01/backend/gunicorn.conf.py --bind 127.0.0.1:8000 --workers 2 --threads 2 backend.wsgi:application
Restart=always
RestartSec=3
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
