# deploy/systemd/bs01-celery-beat.service
[Unit]
Description=BS01 Celery Beat Scheduler
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/BS01/backend
Environment=PYTHONUNBUFFERED=1
Environment=DJANGO_SETTINGS_MODULE=backend.settings
EnvironmentFile=/root/BS01/backend/.env
ExecStart=/root/BS01/.venv/bin/celery -A backend beat -l INFO
# If you install django-celery-beat, you can switch to database scheduler:
# ExecStart=/root/BS01/venv/bin/celery -A backend beat -l INFO --scheduler django_celery_beat.schedulers:DatabaseScheduler
Restart=always
RestartSec=3
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
